# connectors/concur/src/main/service/rhel/files/systemd/concur-connector.service
[Unit]
Description=Concur Connector Service
After=network.target

[Service]
Type=simple
User=roswell
ExecStart=/bin/bash /opt/vmware/connectors/concur/start.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
